# proxy.conf
description "proxy task"

start on runlevel [2345]
stop on runlevel [016]

exec python /home/wendell/python-scripts/proxy/run.py
